# 30_TCC/31_Theory/05_关键技术/local-worker-queue.timer
[Unit]
Description=Run local worker queue dispatcher periodically

[Timer]
OnBootSec=30
OnUnitActiveSec=60
Unit=local-worker-queue.service

[Install]
WantedBy=timers.target
